# deployment/wanda-frontend.service
[Unit]
Description=Wanda Telescope Next.js Frontend
After=network.target wanda-backend.service
Requires=wanda-backend.service
Documentation=https://github.com/AlterMundi/wanda-telescope

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/wanda-telescope/wanda-telescope
Environment="NODE_ENV=production"
Environment="PORT=3000"
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
MemoryMax=1G
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

